# contrib/init/bitcoinreloadedd.conf
description "BitcoinReloaded Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env BITCOINRELOADEDD_BIN="/usr/bin/bitcoinreloadedd"
env BITCOINRELOADEDD_USER="bitcoinreloaded"
env BITCOINRELOADEDD_GROUP="bitcoinreloaded"
env BITCOINRELOADEDD_PIDDIR="/var/run/bitcoinreloadedd"
# upstart can't handle variables constructed with other variables
env BITCOINRELOADEDD_PIDFILE="/var/run/bitcoinreloadedd/bitcoinreloadedd.pid"
env BITCOINRELOADEDD_CONFIGFILE="/etc/bitcoinreloaded/bitcoinreloaded.conf"
env BITCOINRELOADEDD_DATADIR="/var/lib/bitcoinreloadedd"

expect fork

respawn
respawn limit 5 120
kill timeout 600

pre-start script
    # this will catch non-existent config files
    # bitcoinreloadedd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$BITCOINRELOADEDD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run bitcoinreloadedd."
        echo "The setting must appear in $BITCOINRELOADEDD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"BitcoinReloaded Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$BITCOINRELOADEDD_PIDDIR"
    chmod 0755 "$BITCOINRELOADEDD_PIDDIR"
    chown $BITCOINRELOADEDD_USER:$BITCOINRELOADEDD_GROUP "$BITCOINRELOADEDD_PIDDIR"
    chown $BITCOINRELOADEDD_USER:$BITCOINRELOADEDD_GROUP "$BITCOINRELOADEDD_CONFIGFILE"
    chmod 0660 "$BITCOINRELOADEDD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$BITCOINRELOADEDD_PIDFILE" \
    --chuid $BITCOINRELOADEDD_USER:$BITCOINRELOADEDD_GROUP \
    --exec "$BITCOINRELOADEDD_BIN" \
    -- \
    -pid="$BITCOINRELOADEDD_PIDFILE" \
    -conf="$BITCOINRELOADEDD_CONFIGFILE" \
    -datadir="$BITCOINRELOADEDD_DATADIR" \
    -disablewallet \
    -daemon

